# orion-sphere-lrp.service
[Unit]
Description=Orion Sphere LRP Web Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=orion-sphere
Group=orion-sphere
WorkingDirectory=/opt/orion-sphere-lrp
Environment=PATH=/opt/orion-sphere-lrp/venv/bin
ExecStart=/opt/orion-sphere-lrp/venv/bin/python wsgi.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Basic security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orion-sphere-lrp

[Install]
WantedBy=multi-user.target
